# templates/jenkins-agent.service.j2
##############################
#                            #
# Ansible / Tower managed on #
#                            #
##############################

[Unit]
Description=Jenkins Agent: {{ install_jenkins_agent___name }}
Wants=network-online.target
After=network-online.target

[Service]
User={{ install_jenkins_agent___user }}
Group={{ install_jenkins_agent___user }}
Type=simple
Restart=always
WorkingDirectory={{ install_jenkins_agent___home }}
ExecStart=/bin/bash {{ install_jenkins_agent___home }}/jenkins-start-agent.sh

[Install]
WantedBy=multi-user.target
